# scripts/growwiseai-backend-dev.service
[Unit]
Description=GrowWiseAI Backend API (Development)
After=network.target

[Service]
Type=simple
User=sean
WorkingDirectory=/home/sean/GrowWiseAI
Environment="PATH=/home/sean/GrowWiseAI/venv/bin"
ExecStart=/home/sean/GrowWiseAI/venv/bin/uvicorn backend.main:app --reload --host 127.0.0.1 --port 8001
Restart=always
RestartSec=10

# Environment file
EnvironmentFile=/home/sean/GrowWiseAI/googlies.env

# Logging
StandardOutput=append:/var/log/growwiseai-backend-dev.log
StandardError=append:/var/log/growwiseai-backend-dev-error.log

[Install]
WantedBy=multi-user.target
